# etc/periodic/daily/certbot.sh
#!/bin/sh

certbot renew --post-hook "nginx -s reload"
